# systemd/backup.service
# backup.service
[Unit]
Description=Backup home (%h)

[Service]
Type=oneshot
ExecStart=%h/backup_timer.sh
